# sentry/sentry_supervisor.conf
[program:sentry-web]
directory=/srv/sentry/
command=/srv/sentry/bin/sentry --config=/srv/sentry/sentry.conf.py start http
user=sentry
autostart=true
autorestart=true
redirect_stderr=true